# bridges/punchd-bridge-rs/pkg/deb/prerm
#!/bin/sh
set -e

# Stop service before removal
systemctl stop punchd-vpn 2>/dev/null || true
systemctl disable punchd-vpn 2>/dev/null || true
